# vendor/github.com/ligato/vpp-agent/vendor/github.com/ligato/cn-infra/docker/dev_cn_infra/supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log
loglevel=debug
nodaemon=true

[program:agent]
command=/root/go/bin/simple-agent --etcd-config=/opt/agent/dev/etcd.conf --kafka-config=/opt/agent/dev/kafka.conf
autorestart=true
redirect_stderr=true
priority=2
